# etc/speck.conf
description "Speck server"
version "1.0"
author "Joshua Elliott"

respawn
start on startup
chdir /srv/www/speck

exec /srv/www/speck -port=8095 > /srv/www/speck.log

