# etc/speck.conf
description "Speck server"
version "1.0"
author "Joshua Elliott"

respawn
start on startup
chdir /home/josh/speck

exec speck -port=8095 > speck.log

